# illum.service
[Service]
Type=simple
ExecStart=@bindir@/illum-d
Restart=on-failure

[Install]
WantedBy=multi-user.target
